# packaging/systemd/agenc-local-gate-publish@.service
[Unit]
Description=AgenC credential-free local gate publication retry for %i
Documentation=https://github.com/tetsuo-ai/agenc-core/blob/main/docs/ci-required-gates.md
Wants=network-online.target
After=network-online.target
ConditionPathExists=/etc/agenc-local-gatekeeper/config.json

[Service]
Type=oneshot
Slice=system-agencgate.slice
User=root
Group=root
UMask=0077
ExecStart=/usr/bin/flock --no-fork --exclusive /run/agenc-local-gatekeeper/pipeline.lock /opt/agenc-local-gatekeeper/node/bin/node /opt/agenc-local-gatekeeper/repo/scripts/local-gatekeeper.mjs --retry-publish %i
StateDirectory=agenc-local-gatekeeper
StateDirectoryMode=0711
LogsDirectory=agenc-local-gatekeeper
LogsDirectoryMode=0700
RuntimeDirectory=agenc-local-gatekeeper
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
TimeoutStartSec=10min
TimeoutStopSec=30s
KillMode=control-group
SendSIGKILL=yes
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_FOWNER CAP_KILL
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectProc=invisible
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native
TasksMax=128
CPUQuota=100%
MemoryMax=512M
MemorySwapMax=0
OOMPolicy=kill
LimitFSIZE=16M
LimitCORE=0
LimitNOFILE=1024
