# packaging/systemd/agenc-local-gate-dispatcher@.service
[Unit]
Description=AgenC credential-free exact-SHA local gate dispatcher for %i
Documentation=https://github.com/tetsuo-ai/agenc-core/blob/main/docs/ci-required-gates.md
Wants=network-online.target
After=network-online.target
ConditionPathExists=/etc/agenc-local-gatekeeper/config.json

[Service]
Type=oneshot
Slice=system-agencgate.slice
User=root
Group=root
UMask=0077
ExecStart=/usr/bin/flock --no-fork --exclusive /run/agenc-local-gatekeeper/pipeline.lock /opt/agenc-local-gatekeeper/node/bin/node /opt/agenc-local-gatekeeper/repo/scripts/local-gatekeeper.mjs --dispatch %i
StateDirectory=agenc-local-gatekeeper
StateDirectoryMode=0711
LogsDirectory=agenc-local-gatekeeper
LogsDirectoryMode=0700
RuntimeDirectory=agenc-local-gatekeeper
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
TimeoutStartSec=180min
TimeoutStopSec=30s
KillMode=control-group
SendSIGKILL=yes
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_KILL CAP_SETGID CAP_SETUID
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectProc=invisible
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native
TasksMax=4096
CPUQuota=800%
MemoryHigh=12G
MemoryMax=16G
MemorySwapMax=0
OOMPolicy=kill
LimitFSIZE=128M
LimitCORE=0
LimitNOFILE=4096
